# debian/vyatta-cluster.postinst.in
#!/bin/bash

update-rc.d -f heartbeat remove >/dev/null

